# docker/dev/supervisord.d/etcd.ini
[program:etcd]
command = /bin/etcd
priority = 1
autorestart = true
stderr_logfile = /tmp/etcd.err
stdout_logfile = /tmp/etcd.log
stdout_logfile_maxbytes = 50MB
stdout_logfile_backups = 10