# docker/dev/supervisord.d/gochat_logic.ini
[program:gochat_logic]
command = /go/src/gochat/bin/gochat.bin -module logic
priority=3
autorestart = true
stderr_logfile = /tmp/gochat_logic.err
stdout_logfile = /tmp/gochat_logic.log
stdout_logfile_maxbytes = 50MB
stdout_logfile_backups = 10