# docker/dev/supervisord.d/gochat_connect_tcp.ini
[program:gochat_connect_tcp]
command = /go/src/gochat/bin/gochat.bin -module connect_tcp
priority=5
autorestart = true
stderr_logfile = /tmp/gochat_connect_tcp.err
stdout_logfile = /tmp/gochat_connect_tcp.log
stdout_logfile_maxbytes = 50MB
stdout_logfile_backups = 10