# docker/dev/supervisord.d/gochat_api.ini
[program:gochat_api]
command = /go/src/gochat/bin/gochat.bin -module api
priority = 7
autorestart = true
stderr_logfile = /tmp/gochat_api.err
stdout_logfile = /tmp/gochat_api.log
stdout_logfile_maxbytes = 50MB
stdout_logfile_backups = 10