# docker/dev/supervisord.d/gochat_connect.ini
[program:gochat_connect_websocket]
command = /go/src/gochat/bin/gochat.bin -module connect_websocket
priority=4
autorestart = true
stderr_logfile = /tmp/gochat_connect_websocket.err
stdout_logfile = /tmp/gochat_connect_websocket.log
stdout_logfile_maxbytes = 50MB
stdout_logfile_backups = 10