# docker/dev/supervisord.d/redis.ini
[program:redis]
command = /bin/redis-server /root/redis.conf
priority = 2
autorestart = true
stderr_logfile = /tmp/redis.err
stdout_logfile = /tmp/redis.log
stdout_logfile_maxbytes = 50MB
stdout_logfile_backups = 10