# docker/prod/supervisord.d/gochat_site.ini
[program:gochat_site]
command = /go/src/gochat/bin/gochat.bin -module site
priority = 8
autorestart = true
stderr_logfile = /tmp/gochat_site.err
stdout_logfile = /tmp/gochat_site.log
stdout_logfile_maxbytes = 50MB
stdout_logfile_backups = 10