# docker/prod/supervisord.d/gochat_task.ini
[program:gochat_task]
command = /go/src/gochat/bin/gochat.bin -module task
priority = 6
autorestart = true
stderr_logfile = /tmp/gochat_task.err
stdout_logfile = /tmp/gochat_task.log
stdout_logfile_maxbytes = 50MB
stdout_logfile_backups = 10