# docker/dev/etc/supervisord.d/gochat_connect.ini
[program:gochat_connect]
command = /go/src/gochat/bin/gochat_connect.bin
priority=2
autorestart = true
stderr_logfile = /tmp/gochat_connect.err
stdout_logfile = /tmp/gochat_connect.log
stdout_logfile_maxbytes = 50MB
stdout_logfile_backups = 10